# patches/start_screen
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	jdcloud,ax6600)
		ax6600_led >/dev/null 2>&1 &
		;;
	esac
}
